# roles/filebeat/templates/filebeat.service.j2
[Unit]
{%- for filebeat_instance in filebeat_instances %}
Wants={% filebeat_instance %}.filebeat
{% endfor %}
