# roles/filebeat/templates/filebeat.target.j2
[Unit]
Description=Filebeat services group
After=network.target

[Install]
WantedBy=multi-user.target
